# roles/rpi_network/templates/hue-bridge-direct-net.service.j2
[Unit]
Description=Hue bridge direct Ethernet fallback and optional LAN exposure
Wants=network-online.target
After=network-online.target
Documentation=file:{{ install_dir }}/docs/RPI_NETWORK_MODES.md

[Service]
Type=simple
EnvironmentFile=/etc/default/hue-bridge-direct-net
ExecStart=/usr/local/sbin/hue-bridge-direct-net start
ExecStop=/usr/local/sbin/hue-bridge-direct-net stop
TimeoutStartSec=90
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
